# etc/systemd/system/hosting.service
[Unit]
Description=Manage shared hosting accounts

[Service]
Type=simple
ExecStart=/usr/bin/php /var/www/cron.php
