# etc/systemd/system/hosting-del.service
[Unit]
Description=Delete old shared hosting accounts

[Service]
Type=simple
ExecStart=/usr/bin/php /home/daniel/hosting/find_old.php
