# etc/systemd/system/hosting-del.service
[Unit]
Description=Delete old shared hosting accounts

[Service]
Type=simple
ExecStart=/usr/bin/php /var/www/find_old.php
