# etc/systemd/system/tor@.service.d/custom.conf
[Service]
ReadWriteDirectories=-/var/log/tor-instances
